# xcancel-bot.service
[Unit]
Description=XCancel Discord Bot
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/home/deploy/xcant
ExecStart=/home/deploy/.asdf/shims/bundle exec ruby bot.rb
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
